# guest-assets/etc/systemd/system/zfsbox.target
[Unit]
Description=zfsbox minimal boot target
AllowIsolate=yes
